# aniol_gitpull.timer
[Unit]
Description=Activate aniol_gitpull.service every 30 minutes

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s

[Install]
WantedBy=timers.target
