# aniol_gitpull.timer
[Unit]
Description=Activate aniol_gitpull.service every 30 minutes

[Timer]
OnUnitActiveSec=30m
OnBootSec=30m

[Install]
WantedBy=timers.target
